# tests/labs/wireguard/b/lib/systemd/system/netkit-validation.service
[Unit]
Description=Netkit Validation script

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/root/test.sh

[Install]
WantedBy=multi-user.target
