# tests/labs/abr/a/lib/systemd/system/netkit-validation.service
[Unit]
Description=Netkit Validation script

[Service]
Type=simple
ExecStart=/root/test.sh

[Install]
WantedBy=multi-user.target
